# pi/systemd/vivarium-bridge.service
# =============================================================================
# vivarium-bridge.service
# =============================================================================
# Pi bridge — MQTT ↔ Serial ↔ Arduino.
# Starts after provisioner has written device.conf, auto-restarts on failure.
#
# {{SERVICE_USER}} / {{SERVICE_GROUP}} are substituted by setup.sh from the
# account that ran  sudo bash setup.sh  ($SUDO_USER).  Do not hard-code "pi"
# here — modern Raspberry Pi OS lets operators choose their own username.
# =============================================================================

[Unit]
Description=Vivarium Gantry Bridge (Pi MQTT ↔ Serial ↔ Arduino)
After=network-online.target vivarium-provisioner.service
Wants=network-online.target
# Require provisioner to have completed successfully before starting.
Requires=vivarium-provisioner.service

[Service]
Type=simple
User={{SERVICE_USER}}
Group={{SERVICE_GROUP}}
WorkingDirectory={{INSTALL_DIR}}

# Run bridge.py directly — works whether pi/ is on sys.path or not.
ExecStart=/usr/bin/python3 {{INSTALL_DIR}}/bridge.py

# Auto-restart on failure, 5s cooldown.
# After 10 failures in 60s, systemd marks it failed so an admin can investigate.
Restart=on-failure
RestartSec=5s
StartLimitBurst=10
StartLimitIntervalSec=60

# All output goes to the journal: journalctl -u vivarium-bridge -f
StandardOutput=journal
StandardError=journal

Environment=PYTHONUNBUFFERED=1
Environment=GANTRY_LOG_LEVEL=INFO
# Override without editing this file:
#   sudo systemctl edit vivarium-bridge
#   [Service]
#   Environment=GANTRY_LOG_LEVEL=DEBUG

NoNewPrivileges=true

[Install]
WantedBy=multi-user.target
