# pi/systemd/vivarium-camera.service
# =============================================================================
# vivarium-camera.service
# =============================================================================
# MediaMTX live-streaming agent for this Pi.
#
# Relationship to other units:
#   • INDEPENDENT of vivarium-bridge.service — if the bridge crashes the
#     stream keeps running, and vice versa (Section 5.6).
#   • After= vivarium-provisioner.service so mediamtx.yaml exists before start.
#   • Does NOT Require= the bridge — stream independence is intentional.
#
# MediaMTX binary: /usr/local/bin/mediamtx  (installed by setup.sh)
# Config file:     {{INSTALL_DIR}}/mediamtx/mediamtx.yaml  (written by provisioner)
#
# {{SERVICE_USER}} / {{SERVICE_GROUP}} are substituted by setup.sh from the
# account that ran  sudo bash setup.sh  ($SUDO_USER).  Do not hard-code "pi"
# here — modern Raspberry Pi OS lets operators choose their own username.
#
# Ports:
#   :8554  — RTSP   (Pi camera → server pull)
#   :8888  — HTTP / MJPEG
#   :8889  — WebRTC (WHEP)
#   127.0.0.1:9997  — REST API (health checks, localhost only)
#
# Logging: journalctl -u vivarium-camera -f
# =============================================================================

[Unit]
Description=Vivarium Gantry — MediaMTX live-streaming agent
After=network-online.target vivarium-provisioner.service
Wants=network-online.target
# Independent of the bridge — intentional.

[Service]
Type=simple
User={{SERVICE_USER}}
Group={{SERVICE_GROUP}}
WorkingDirectory={{INSTALL_DIR}}/mediamtx

# MediaMTX binary installed by setup.sh to /usr/local/bin/mediamtx
ExecStart=/usr/local/bin/mediamtx {{INSTALL_DIR}}/mediamtx/mediamtx.yaml

# Self-heal on failure — 5s cooldown, max 10 restarts per minute.
Restart=on-failure
RestartSec=5s
StartLimitBurst=10
StartLimitIntervalSec=60

StandardOutput=journal
StandardError=journal
SyslogIdentifier=vivarium-camera

# Pi Camera requires access to /dev/video* and /dev/dma_heap.
# {{SERVICE_USER}} must be in the "video" group — setup.sh adds them if needed.
SupplementaryGroups=video

# Security hardening
NoNewPrivileges=true

# MediaMTX tuning — limit to 1 core, aggressive GC to save RAM on Pi
Environment=GOGC=50
Environment=GOMAXPROCS=1

[Install]
WantedBy=multi-user.target
