# pi/systemd/vivarium-camera.service
# =============================================================================
# vivarium-camera.service  (Section 5.1 / 5.6 / 8)
# =============================================================================
#
# go2rtc live-streaming agent for this Pi.
#
# Relationship to other units:
#   • Completely INDEPENDENT of vivarium-bridge.service — if the bridge crashes,
#     the stream keeps running, and vice versa (Section 5.6).
#   • After= vivarium-provisioner.service so it starts only after device.conf
#     (and therefore /etc/gantry/go2rtc.yaml) exist.
#   • Does NOT Require= or BindsTo= the bridge — stream independence is
#     intentional and must not be weakened.
#
# Config file: /etc/gantry/go2rtc.yaml
#   Written by provisioner.py using pi/go2rtc/go2rtc.example.yaml as the
#   template (Section 5.3). Contains the device_id (stream name), RTSP port,
#   and rtsp_password issued by POST /provision.
#
# Ports:
#   :{{GO2RTC_PORT}} (default 8554)  — RTSP, externally reachable by the server
#   127.0.0.1:1984                   — HTTP API, localhost only
#
# Logging: journalctl -u vivarium-camera [-f]
# =============================================================================

[Unit]
Description=Vivarium Gantry — go2rtc live-streaming agent (Section 5.6 / 8)
Documentation=https://github.com/AlexxIT/go2rtc
# Intentionally does NOT Require= vivarium-bridge.service — stream independence
# is a hard requirement from Section 5.6.
After=network-online.target vivarium-provisioner.service
Wants=network-online.target

[Service]
Type=simple
User=gantry
Group=gantry
WorkingDirectory=/opt/vivarium/pi/go2rtc

# go2rtc binary must be installed at /usr/local/bin/go2rtc.
# Install it during the OS image build or via the Ansible playbook.
# The config is generated by provisioner.py from go2rtc.example.yaml.
ExecStart=/usr/local/bin/go2rtc -config /etc/gantry/go2rtc.yaml

# Restart immediately on failure so the stream self-heals without operator
# intervention. 5s delay avoids a tight crash loop on config errors.
Restart=on-failure
RestartSec=5s
# After 10 fast failures in 60s, systemd marks the unit as failed so
# an admin can investigate rather than spinning forever.
StartLimitBurst=10
StartLimitIntervalSec=60

# Structured logging — all output goes to the journal.
StandardOutput=journal
StandardError=journal
SyslogIdentifier=vivarium-camera

# ── Security hardening (Section 9 Layer 2D) ───────────────────────────────
# Items marked [PROD ONLY] are inert unless their prerequisite is set up.
NoNewPrivileges=true
# go2rtc only needs read access to its config file; it holds streams in memory.
ReadOnlyPaths=/etc/gantry/go2rtc.yaml
# [PROD ONLY] Uncomment when the filesystem is hardened:
# ProtectSystem=strict
# ProtectHome=true
# PrivateTmp=true

# go2rtc process tuning — single core, aggressive GC to limit RAM on a Pi.
Environment=GOGC=50
Environment=GOMAXPROCS=1
# [PROD ONLY] Override log level without editing this file:
#   sudo systemctl edit vivarium-camera
#   [Service]
#   Environment=GO2RTC_LOG=debug

[Install]
WantedBy=multi-user.target
